# docker/frankenphp/s6/frankenphp/run
#!/command/with-contenv sh
# Load the secrets generated by init.sh (.env.local) into the process
# environment — Caddy reads {env.MERCURE_JWT_SECRET} at startup.
if [ -f /var/www/html/.env.local ]; then
  set -a
  . /var/www/html/.env.local
  set +a
fi
exec frankenphp run --config /etc/caddy/Caddyfile
